# contrib/init/dogecoind.conf
description "Botccoin Core Daemon"

start on runlevel [2345]
stop on starting rc RUNLEVEL=[016]

env BotccoinD_BIN="/usr/bin/Botccoind"
env BotccoinD_USER="Botccoin"
env BotccoinD_GROUP="Botccoin"
env BotccoinD_PIDDIR="/var/run/Botccoind"
# upstart can't handle variables constructed with other variables
env BotccoinD_PIDFILE="/var/run/Botccoind/Botccoind.pid"
env BotccoinD_CONFIGFILE="/etc/Botccoin/Botccoin.conf"
env BotccoinD_DATADIR="/var/lib/Botccoind"

expect fork

respawn
respawn limit 5 120
kill timeout 60

pre-start script
    # this will catch non-existent config files
    # Botccoind will check and exit with this very warning, but it can do so
    # long after forking, leaving upstart to think everything started fine.
    # since this is a commonly encountered case on install, just check and
    # warn here.
    if ! grep -qs '^rpcpassword=' "$BotccoinD_CONFIGFILE" ; then
        echo "ERROR: You must set a secure rpcpassword to run Botccoind."
        echo "The setting must appear in $BotccoinD_CONFIGFILE"
        echo
        echo "This password is security critical to securing wallets "
        echo "and must not be the same as the rpcuser setting."
        echo "You can generate a suitable random password using the following"
        echo "command from the shell:"
        echo
        echo "bash -c 'tr -dc a-zA-Z0-9 < /dev/urandom | head -c32 && echo'"
        echo
        echo "It is also recommended that you also set alertnotify so you are "
        echo "notified of problems:"
        echo
        echo "ie: alertnotify=echo %%s | mail -s \"Botccoin Alert\"" \
            "admin@foo.com"
        echo
        exit 1
    fi

    mkdir -p "$BotccoinD_PIDDIR"
    chmod 0755 "$BotccoinD_PIDDIR"
    chown $BotccoinD_USER:$BotccoinD_GROUP "$BotccoinD_PIDDIR"
    chown $BotccoinD_USER:$BotccoinD_GROUP "$BotccoinD_CONFIGFILE"
    chmod 0660 "$BotccoinD_CONFIGFILE"
end script

exec start-stop-daemon \
    --start \
    --pidfile "$BotccoinD_PIDFILE" \
    --chuid $BotccoinD_USER:$BotccoinD_GROUP \
    --exec "$BotccoinD_BIN" \
    -- \
    -pid="$BotccoinD_PIDFILE" \
    -conf="$BotccoinD_CONFIGFILE" \
    -datadir="$BotccoinD_DATADIR" \
    -disablewallet \
    -daemon

